# animebg.timer
[Unit]
Description=animebg scheduler

[Timer]
OnBootSec=0min
OnUnitActiveSec=1min

[Install]
WantedBy=graphical.target